# networking/smf/manta-nic.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
    This Source Code Form is subject to the terms of the Mozilla Public
    License, v. 2.0. If a copy of the MPL was not distributed with this
    file, You can obtain one at http://mozilla.org/MPL/2.0/.
-->

<!--
    Copyright (c) 2014, Joyent, Inc.
-->

<service_bundle type='manifest' name='export'>
  <service name='smartdc/hack/manta-nic' type='service' version='0'>

    <create_default_instance enabled='true'/>
    <single_instance/>

    <dependency name='fs-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    <dependency name='network' grouping='require_all' restart_on='restart' type='service'>
      <service_fmri value='svc:/milestone/network'/>
    </dependency>

    <exec_method name='start' type='method' exec='/opt/custom/mnet/bin/manta-nic.sh' timeout_seconds='60'>
      <method_context>
        <method_credential user="root" group="root" />
      </method_context>
    </exec_method>

    <exec_method name='stop' type='method' exec=':true' timeout_seconds='60'/>

    <property_group name='startd' type='framework'>
      <propval name='duration' type='astring' value='transient' />
    </property_group>

    <template>
      <common_name>
        <loctext xml:lang='C'>A hack for creating a nic on the manta network</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
